# deploy/systemd/digitmak-certbot-renew.service
[Unit]
Description=Renew DigitMak Portal TLS certificates
Requires=docker.service
After=docker.service digitmak-portal.service

[Service]
Type=oneshot
WorkingDirectory=/opt/digitmak-portal/current
ExecStart=/usr/bin/docker compose -p digitmak-portal --env-file /opt/digitmak-portal/shared/.env.production -f docker-compose.production.yml run --rm certbot renew --webroot -w /var/www/certbot --quiet
ExecStartPost=/usr/bin/docker compose -p digitmak-portal --env-file /opt/digitmak-portal/shared/.env.production -f docker-compose.production.yml exec -T nginx nginx -s reload
